# yggdrasil-iptunnel/files/etc/systemd/system/yggdrasil.service.d/yggdrasil-iptunnel.conf
[Service]
ExecStartPost=/usr/sbin/yggdrasil-setup
