# ssb/files/DEBIAN/postinst
#!/bin/bash

adduser --system ssb-server

systemctl daemon-reload
systemctl enable ssb.service
systemctl start ssb.service
