# ssb/files/DEBIAN/prerm
systemctl stop ssb.service
systemctl disable ssb.service
systemctl daemon-reload
