# prometheus/files/DEBIAN/postinst
#!/bin/bash
systemctl enable babeld
systemctl start babeld
